# overlay.d/99my/usr/lib/systemd/system/etcd-ready.target
[Unit]
Description=Etcd is Running and available for Kubernetes
Requires=multi-user.target
After=etcd.service

[Install]
WantedBy=multi-user.target
