# overlay.d/99my/usr/lib/systemd/system/etcd.path
[Unit]
Description=Etcd service can be activated if cert and config is given

[Path]
PathExists=/etc/etcd/ca.pem /etc/etcd/kubernetes.pem /etc/etcd/kubernetes-key.pem
Unit=k8sos-etcd.service
