# overlay.d/99my/usr/lib/systemd/system/etcd-ready.target
[Unit]
Description=Etcd is Running and available for Kubernetes
Requires=multi-user.target k8sos-etcd.service
After=k8sos-etcd.service
